# .config/etc/crontab
*/2 * * * * /home/rj1/.config/email/check-email.sh
