# .config/etc/sv/user-services/run
#!/bin/sh

export USER="rj1"
export HOME="/home/$USER"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.config/runit/services"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
