# .config/runit/services/miniircd/run
#!/bin/sh

exec /home/rj1/bin/miniircd 4>&1
